# supervisord.conf
[supervisord]
user=root

[program:docker]
user=root
autostart=true
autorestart=true
command=/usr/bin/docker -d
